# docker/ssh/services.d/openssh/run
#!/usr/bin/execlineb -P
/usr/sbin/sshd -p 22 -h /etc/ssh/ssh_host_dsa_key