# docker/php/services.d/caddy/run
#!/usr/bin/execlineb -P
s6-setuidgid www-data
caddy run --config /var/www/conf/Caddyfile